# systemd/scheduLight.target
[Unit]
Description=scheduLight
[Install]
WantedBy=multi-user.target
